# .scripts/sync-repos.upstart.conf
start on runlevel [06]

exec --chuid minhajuddin --exec  /bin/bash -l -c /home/minhajuddin/.scripts/sync-repos
